# root/etc/uci-defaults/50_luci-app-autoupdate
#!/bin/sh

# 确保 RPC 脚本有执行权限
chmod 755 /usr/libexec/rpcd/autoupdate.uc

# 注册修复脚本
chmod 755 /usr/bin/autoupdate-fix.sh

# 添加防火墙允许规则
uci -q batch <<-END >/dev/null
	set firewall.autoupdate=rule
	set firewall.autoupdate.name='Allow-Autoupdate'
	set firewall.autoupdate.src='wan'
	set firewall.autoupdate.target='ACCEPT'
	set firewall.autoupdate.enabled='1'
	commit firewall
END

# 重启rpcd服务以加载新的RPC方法
/etc/init.d/rpcd stop
sleep 1
/etc/init.d/rpcd start

# 如果服务没有自动注册，尝试手动注册
if ! ubus list | grep -q autoupdate; then
    /usr/bin/autoupdate-fix.sh
fi

# 移除LuCI缓存
rm -f /tmp/luci-indexcache
exit 0
	set firewall.autoupdate.name='Allow-Autoupdate'
	set firewall.autoupdate.src='wan'
	set firewall.autoupdate.target='ACCEPT'
	set firewall.autoupdate.enabled='1'
	commit firewall
END

# 立即重启rpcd服务
/etc/init.d/rpcd stop
sleep 1
/etc/init.d/rpcd start

# 清除LuCI缓存
rm -f /tmp/luci-indexcache
exit 0
